# mk/log_tty.conf

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn
exec /usr/bin/tail -F /var/log/messages > /dev/tty4 2>&1 < /dev/null
