# startup/60-firmware.rules
SUBSYSTEM=="firmware", RUN+="firmware.agent"
